# cron/nyc-eats-refresh.service.in
[Unit]
Description=NYC Eats daily data refresh & cross-reference
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User={{VPS_USER}}
Environment=NYC_EATS_REPO={{VPS_REPO}}
Environment=NYC_EATS_SERVE={{VPS_PATH}}
Environment=NYC_EATS_LOG={{VPS_PATH}}/refresh.log
Environment=YELP_API_KEY={{YELP_API_KEY}}
Environment=GOOGLE_API_KEY={{GOOGLE_API_KEY}}
ExecStart={{VPS_REPO}}/cron/nyc-eats-refresh
StandardOutput=append:{{VPS_PATH}}/refresh.log
StandardError=append:{{VPS_PATH}}/refresh.log

# Hardening
NoNewPrivileges=true
PrivateTmp=true
